# scripts/rpm/prerm
#!/bin/sh

# Abort if any command returns an error value
set -e

if [ -f "/etc/init.d/centengine" ]; then
    update-rc.d centengine remove
fi

exit 0
